# profile.d/java.sh
# DONT CHANGE
export PATH=/opt/play:$PATH
export PATH=/opt/intellij/bin:$PATH
