# profile.d/scala.sh
# DONT CHANGE
export PATH=/opt/scalaide:$PATH
